# web/systemd/octollama-web.service
[Unit]
Description=OctoLlama - panel WWW
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/octollama-web
ExecStart=%h/.local/share/octollama-web/.venv/bin/waitress-serve --host=0.0.0.0 --port={{PORT}} app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
