# daemon/systemd/octollama-daemon.service
[Unit]
Description=OctoLlama - lokalny demon sterujacy usluga Ollama
After=network-online.target remote-fs.target
Wants=network-online.target
RequiresMountsFor=/var/lib/octollama/state

[Service]
Type=simple
ExecStart=/opt/octollama/daemon/.venv/bin/python3 /opt/octollama/daemon/octollama_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
