# daemon/systemd/ollama-manager-daemon.service
[Unit]
Description=Ollama Manager - lokalny demon sterujacy usluga Ollama
After=network.target

[Service]
Type=simple
ExecStart=/opt/ollama-manager/daemon/.venv/bin/python3 /opt/ollama-manager/daemon/ollama_manager_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
